# infra/algorhythm.service
[Unit]
Description=AlgoRhythm FastAPI inference server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/algorhythm/backend
# '-' = don't fail to start if the file is missing yet
EnvironmentFile=-/opt/algorhythm/backend/.env
# 1 worker — t3.micro is 1GB; a 2nd worker loading the ML stack would OOM
ExecStart=/opt/algorhythm/backend/.venv/bin/gunicorn main:app \
    --workers 1 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --timeout 120
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
